# provisioning/roles/java7/files/etc/profile.d/maven.sh
export PATH="/opt/maven/current/bin:${PATH}"
