# provisioning/roles/spark/files/etc/profile.d/spark.sh
#! /bin/bash

export YARN_CONF_DIR=/etc/hadoop/conf
